# news/leafnode/files/org.macports.leafnode.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>OnDemand</key>			<true/>
		<key>Label</key>			<string>org.macports.leafnode</string>
		<key>WorkingDirectory</key>	<string>@PREFIX@/var/spool/news</string>
		<key>Program</key>			<string>/usr/bin/sudo</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/bin/sudo</string>
			<string>-u</string>
			<string>_news</string>
			<string>@PREFIX@/sbin/leafnode</string>
		</array>
		<key>Sockets</key>
		<dict>
			<key>Listeners</key>
			<dict>
				<key>SockServiceName</key>	<string>nntp</string>
			</dict>
		</dict>
		<key>inetdCompatibility</key>
		<dict>
			<key>Wait</key>	<false/>
		</dict>
	</dict>
</plist>

<!--vim: set nowrap tabstop=4 shiftwidth=4 softtabstop=4 noexpandtab :-->
<!--vim: set textwidth=0 filetype=xml foldmethod=marker nospell :-->
