# news/leafnode/files/org.macports.texpire.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>KeepAlive</key>		<false/>
		<key>Label</key>			<string>org.macports.texpire</string>
		<key>StartInterval</key>	<integer>86400</integer>
		<key>WorkingDirectory</key>	<string>@PREFIX@/var/spool/news</string>
		<key>Program</key>			<string>/usr/bin/sudo</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/bin/sudo</string>
			<string>-u</string>
			<string>_news</string>
			<string>@PREFIX@/sbin/texpire</string>
		</array>
	</dict>
</plist>

<!--vim: set nowrap tabstop=4 shiftwidth=4 softtabstop=4 noexpandtab :-->
<!--vim: set textwidth=0 filetype=xml foldmethod=marker nospell :-->
