# news/leafnode/files/org.macports.fetchnews.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>org.macports.fetchnews</string>
	<key>Program</key>
	<string>/usr/bin/sudo</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/sudo</string>
		<string>-u</string>
		<string>_news</string>
		<string>@PREFIX@/sbin/fetchnews</string>
	</array>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>WorkingDirectory</key>
	<string>@PREFIX@/var/spool/news</string>
</dict>
</plist>
