# security/gpg-agent/files/gpg-agent.plist.default
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>@LABEL@</string>

	<!-- Please uncomment on 10.4; OnDemand doesn't work properly there. -->
	<!--
	<key>OnDemand</key>
	<false />
	-->

	<key>ProgramArguments</key>
	<array>
		<string>@PREFIX@/bin/gpg-agent</string>
		<string>--write-env-file</string>
		<string>--no-use-standard-socket</string>
		<string>--launchd</string>
		<!--@SSH_SUPPORT@
		<string>--enable-ssh-support</string>
		@SSH_SUPPORT@-->
	</array>

	<key>RunAtLoad</key>
	<true/>

	<!-- Restart the job, should it crash. -->
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false />
	</dict>

	<key>ServiceIPC</key>
	<true/>

	<key>Sockets</key>
	<dict>
		<key>Listeners_agent</key>
		<dict>
			<key>SecureSocketWithKey</key>
			<string>GPG_AGENT_INFO</string>
		</dict>

		<!--@SSH_SUPPORT@
		<key>Listeners_ssh</key>
		<dict>
			<key>SecureSocketWithKey</key>
			<string>SSH_AUTH_SOCK</string>
		</dict>
		@SSH_SUPPORT@-->
	</dict>
</dict>
</plist>
