# devel/rtags/files/rdm.plist.in
<?xml version='1.0' encoding='UTF-8'?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd" >
<plist version='1.0'>
    <dict>
        <key>Label</key>
        <string>@LABEL@</string>

        <key>ProgramArguments</key>
        <array>
            <string>@PREFIX@/bin/rdm</string>
            <string>-v</string>
            <string>--launchd</string>
            <string>--socket-file</string>
            <string>@RDM_SOCKET_PATH@</string>
            <string>--inactivity-timeout</string>
            <string>300</string>
        </array>

        <key>Disabled</key>
        <true/>

        <key>Sockets</key>
        <dict>
            <key>Listener</key>
            <dict>
                <key>SockPathName</key>
                <string>@RDM_SOCKET_PATH@</string>
                <key>SockPathMode</key>
                <integer>438</integer> <!-- 0666 -->
            </dict>
        </dict>
    </dict>
</plist>
